# ops/systemd/user/logrotate-m5-trader.service
[Unit]
Description=Logrotate for m5_trader logs

[Service]
Type=oneshot
ExecStart=%h/bin/logrotate_m5_trader.sh

