# ops/systemd/user/m5-ibg-logs-link.service
[Unit]
Description=Update File-Window ibg_logs symlink to latest IBG logs

[Service]
Type=oneshot
ExecStart=%h/File-Window/tools/update_ibg_log_symlink.sh
# Treat exit code 1 (no logs yet) as success to avoid spam
SuccessExitStatus=1
