# ops/systemd/user/m5trader-nightly.service
[Unit]
Description=File-Window nightly validation

[Service]
Type=oneshot
ExecStart=%h/File-Window/scripts/nightly_validate.sh
