# ops/systemd/user/m5trader-watch.service
[Unit]
Description=File-Window autosync watcher
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/File-Window
ExecStart=%h/File-Window/bin/m5-watch.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
