# ops/systemd/user/m5trader-scan.service
[Unit]
Description=File-Window periodic autosync

[Service]
Type=oneshot
WorkingDirectory=%h/File-Window
ExecStart=%h/File-Window/bin/m5-sync.sh --if-needed
