# ops/systemd/user/m5trader-nightly.service
[Unit]
Description=M5-Trader nightly validation

[Service]
Type=oneshot
ExecStart=%h/M5-Trader/scripts/nightly_validate.sh

