# ops/systemd/user/m5trader-nightly.timer
[Unit]
Description=File-Window nightly validation timer

[Timer]
OnCalendar=02:30
Persistent=true
Unit=m5trader-nightly.service

[Install]
WantedBy=timers.target
